# com.fitlocal.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fitlocal.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brendanreed/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>packages/api/dist/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brendanreed/Projects/fitlocal</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3001</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/brendanreed/Library/Logs/fitlocal.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brendanreed/Library/Logs/fitlocal.log</string>
</dict>
</plist>
